# assets/files/resumator.conf
[program:resumator]
command=/opt/resumator/resumator.linux -c /opt/resumator/resumator.conf
redirect_stderr=true
stdout_logfile=/var/log/resumator.log
stdout_logfile_maxbytes=100MB
